# scripts/com.caiso.ml_pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Job identifier - must be unique -->
    <key>Label</key>
    <string>com.caiso.ml_pipeline</string>
    
    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/bran/Work/model-power-nowcast/.venv/bin/python</string>
        <string>/Users/bran/Work/model-power-nowcast/scripts/automated_ml_pipeline.py</string>
        <string>--target-zones</string>
        <string>NP15</string>
        <string>--quiet</string>
    </array>
    
    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/bran/Work/model-power-nowcast</string>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/bran/Work/model-power-nowcast/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/bran/Work/model-power-nowcast</string>
    </dict>
    
    <!-- Schedule: every 6 hours at 2 AM, 8 AM, 2 PM, 8 PM -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>2</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>14</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    
    <!-- Run at load (when system starts) -->
    <key>RunAtLoad</key>
    <false/>
    
    <!-- Keep job alive if it exits -->
    <key>KeepAlive</key>
    <false/>
    
    <!-- Standard output and error logs -->
    <key>StandardOutPath</key>
    <string>/Users/bran/Work/model-power-nowcast/logs/ml_pipeline_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/bran/Work/model-power-nowcast/logs/ml_pipeline_stderr.log</string>
    
    <!-- Process type -->
    <key>ProcessType</key>
    <string>Background</string>
    
    <!-- Nice level (lower priority) -->
    <key>Nice</key>
    <integer>15</integer>
    
    <!-- Throttle interval (minimum time between launches) -->
    <key>ThrottleInterval</key>
    <integer>3600</integer>
</dict>
</plist>
